# deploy/clash-autopull.service
[Unit]
Description=Clash Auto-Deploy Watcher
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/clash
ExecStart=/bin/bash /opt/clash/deploy/autopull.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
